# portal/deploy/systemd/labs-lab1.service
[Unit]
Description=Streamlit lab 1 (численные методы)
After=network.target

[Service]
Type=simple
User=labs
Group=labs
WorkingDirectory=/opt/labs/numerical_labs/lab1
ExecStart=/opt/labs/numerical_labs/lab1/venv/bin/streamlit run /opt/labs/numerical_labs/lab1/homework_ui.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
